# metrics.initscript
#!/bin/sh /etc/rc.common

USE_PROCD=1

start_service() {
  procd_open_instance
  procd_set_param command /bin/metrics.lua -p 9100
  procd_close_instance
}
